# release/ladon.service
[Unit]
Description=ladon — smart split routing for DPI-censored networks
Documentation=https://github.com/belotserkovtsev/ladon
After=network-online.target dnsmasq.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ladon
ExecStart=/opt/ladon/ladon \
  -db /opt/ladon/state/engine.db \
  run \
  -manual-allow /etc/ladon/manual-allow.txt \
  -manual-deny  /etc/ladon/manual-deny.txt \
  /var/log/dnsmasq.log
Restart=always
RestartSec=5

# Hardening — engine only needs CAP_NET_ADMIN for ipset mutations if it ever
# moves from exec'ing the CLI to netlink. Today it's root because `ipset add`
# is root-only anyway.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true
# ladon writes /etc/dnsmasq.d/ladon-manual.conf at startup with ipset=
# directives for manual-allow + extensions, then SIGHUPs dnsmasq.
ReadWritePaths=/etc/dnsmasq.d

[Install]
WantedBy=multi-user.target
